# wifi/fix-wifi.service
[Unit]
Description=Fix Wifi

[Service]
ExecStart=/usr/bin/fix-wifi

[Install]
WantedBy=multi-user.target
